# airootfs/etc/systemd/system/optionstimeout.timer
[Unit]
Description=Run

[Timer]
OnBootSec=1min
Unit=optionstimeout.service

[Install]
WantedBy=multi-user.target
