# airootfs/etc/systemd/system/optionstimeout.service
[Unit]
Description=MyScript

[Service]
Type=simple
ExecStart=/usr/local/bin/ReturnOptions.sh
